# deploy/minions.service
[Unit]
Description=Minions of Darkness web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/minions
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m minions.app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
